# udev-rules/83.ac_power.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/share/slimbookbattery/src/ac_power.py"

