# debian/postinst
#!/bin/sh
# postinst script for slimbookbattery

set -e

python3 /usr/share/slimbookbattery/src/check_config.py

chmod -R 777 /usr/share/slimbookbattery/src/*.py
chmod -R 777 /usr/share/slimbookbattery/src/*.sh

# Enabling tlp.service
sudo systemctl enable tlp.service

#Autostar file changed
FILE=/home/*/.config/autostart/slimbookbattery-autostart.desktop

if test -f $FILE; then
    rm /home/*/.config/autostart/slimbookbattery-autostart.desktop
fi

FILE=/lib/systemd/system/slimbookbattery.service

if test -f $FILE; then
    rm $FILE
fi

exit 0
